# scripts/com.oink.crypto-quant-okx.dashboard-keepalive.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.oink.crypto-quant-okx.dashboard-keepalive</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Volumes/MacHD/Projects/crypto-quant-okx/scripts/keep_dashboard_alive.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>StartInterval</key>
    <integer>1800</integer>

    <key>WorkingDirectory</key>
    <string>/Volumes/MacHD/Projects/crypto-quant-okx</string>

    <key>StandardOutPath</key>
    <string>/Volumes/MacHD/Projects/crypto-quant-okx/logs/dashboard-keepalive.launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/Volumes/MacHD/Projects/crypto-quant-okx/logs/dashboard-keepalive.launchd.err.log</string>
</dict>
</plist>
